# crontab
#minute-hour-day-month-weekday command
0 * * * * /Users/MorehouseJ09/Documents/programs/cron-jobs/sync.sh

