# crontab
#minute-hour-day-month-weekday command
0 12 * * * /Users/MorehouseJ09/Documents/programs/cron-jobs/sync.sh

# 
0 12 * * * /Users/MorehouseJ09/Documents/programs/cron-jobs/clean_history.sh

